# configs/releng/airootfs/etc/systemd/system/x11.service
# x11.service

[Unit]
Description=X11 Window System
After=osjs.service

[Service]
User=osjs
Group=osjs
Restart=on-failure
ExecStart=/usr/bin/startx
